# deploy/systemd/finance-app.service
# Personal Finance App - Gunicorn Service for OpenLiteSpeed/CyberPanel
# =======================================================================
# Copy to: /etc/systemd/system/finance-app.service
# Domain: finance.mstatilitechnologies.com
# CyberPanel User: finan1713
# =======================================================================

[Unit]
Description=Personal Finance App - Gunicorn
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify
# CyberPanel assigned user for finance.mstatilitechnologies.com
User=finan6751
Group=finan6751

# Working directory - where manage.py is located (public_html is git root)
WorkingDirectory=/home/finance.mstatilitechnologies.com/public_html

# Environment file location (outside public_html for security)
EnvironmentFile=/home/finance.mstatilitechnologies.com/.env

# Gunicorn command
ExecStart=/home/finance.mstatilitechnologies.com/.venv/bin/gunicorn \
    --bind 127.0.0.1:8001 \
    --workers 3 \
    --timeout 120 \
    --access-logfile /home/finance.mstatilitechnologies.com/logs/gunicorn-access.log \
    --error-logfile /home/finance.mstatilitechnologies.com/logs/gunicorn-error.log \
    --capture-output \
    backend.wsgi:application

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=5min
StartLimitBurst=10

# Security hardening
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
